# mapper-devusb.service.in
[Unit]
Description=mapper-devusb daemon

[Service]
Type=simple
User=mapper-devusb
ExecStart=@bindir@/mapper-devusb

[Install]
WantedBy=multi-user.target

